# scripts/com.polymarket.heal-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polymarket.heal-loop</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/macbookair/polymarketv2/scripts/heal_loop.sh</string>
    </array>

    <!-- Run continuously; relaunch if it crashes. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Don't hammer launchd if heal_loop somehow exits in <10s -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/polymarket-heal.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/polymarket-heal.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/macbookair/polymarketv2</string>
</dict>
</plist>
